# upstart/lumberjack.conf
start on runlevel [2345]

start on (started network-interface
          or started network-manager
          or started networking)

stop on (stopping network-interface
         or stopping network-manager
         or stopping networking)

respawn

script
  reset_lumberjack=false
  if [ -f "/tmp/reset_lumberjack" ]; then
    reset_lumberjack=true
    rm -f /tmp/reset_lumberjack
  fi
	exec su - deploy -c "cd ~/lumberjack/ ; /opt/logstash-forwarder/bin/logstash-forwarder -config lumberjack.conf -from-beginning=$reset_lumberjack > logstash-forwarder.log"
end script

stop on runlevel [016]
